# private_dot_config/systemd/user/bluetooth-off.service
[Unit]
Description=Power off Bluetooth by default

[Service]
Type=oneshot
ExecStart=/usr/bin/bluetoothctl power off

[Install]
WantedBy=default.target
