# private_dot_config/systemd/user/swaync.service
[Unit]
Description=swaync

[Service]
ExecStart=swaync
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
