# private_dot_config/systemd/user/swayidle.service
[Unit]
Description=swayidle

[Service]
ExecStart=swayidle -w
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
